# tierflow.service
[Unit]
Description=Tierflow storage balancer daemon
Documentation=https://github.com/leonidbkh/tierflow
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Run as user that has access to your storage tiers
# IMPORTANT: Change this to your user or create dedicated user
User=root

# Command to run
ExecStart=/usr/local/bin/tierflow daemon --config /etc/tierflow/config.yaml --interval 3600

# Restart policy
Restart=on-failure
RestartSec=60

# Timeouts
TimeoutStopSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tierflow

# Resource limits - run with low priority
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

# Environment
# Uncomment to enable debug logging:
# Environment="RUST_LOG=tierflow=debug"

[Install]
WantedBy=multi-user.target
